# var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.rFxbRw/crontab installed on Tue Dec 29 09:43:51 2015)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
*/5 * * * * /home/pi/shell/Python/FileUpload_01.py > /dev/null 2>&1
*/5 * * * * /home/pi/shell/Python/FileUpload_02.py > /dev/null 2>&1
*/5 * * * * /home/pi/shell/Python/FileUpload_03.py > /dev/null 2>&1
*/5 * * * * /home/pi/shell/Python/FileUpload_04.py > /dev/null 2>&1
